# etc/periodic/daily/000.tarsnap-backup
#!/bin/sh

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

paths=$daily_tarsnap_backup_paths
if [ -z "$paths" ]; then
    paths='/'
fi

keep=$daily_tarsnap_backup_keep
if [ -z "$keep" ]; then
    keep=7
fi

case "$daily_tarsnap_backup_enable" in
    [Yy][Ee][Ss])
        . /bin/tarsnap-backup
        do_backups "$paths" $keep 'daily'
        ;;
    *)
        ;;
esac
