# etc/periodic/weekly/000.tarsnap-backup
#!/bin/sh

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

paths=$weekly_tarsnap_backup_paths
if [ -z "$paths" ]; then
    paths='/'
fi

keep=$weekly_tarsnap_backup_keep
if [ -z "$keep" ]; then
    keep=7
fi

case "$weekly_tarsnap_backup_enable" in
    [Yy][Ee][Ss])
        tarsnap-backup "$paths" $keep 'weekly'
        ;;
    *)
        ;;
esac
